# unit-files/example-com-server.service
[Service]
Type=forking
Environment=NAMEDCONF=/vagrant/example-com-server/named.conf
ExecStart=/usr/sbin/named -u named -c ${NAMEDCONF} 


