# unit-files/notsigned-com-server.service
[Service]
Type=forking
Environment=NAMEDCONF=/vagrant/notsigned-com-server/named.conf
ExecStart=/usr/sbin/named -u named -c ${NAMEDCONF} 


