# unit-files/resolver.service
[Service]
Type=forking
Environment=NAMEDCONF=/vagrant/resolver/named.conf
ExecStart=/usr/sbin/named -u named -c ${NAMEDCONF} 
