# unit-files/root-server.service
[Service]
Type=forking
Environment=NAMEDCONF=/vagrant/root-server/named.conf
ExecStart=/usr/sbin/named -u named -c ${NAMEDCONF} 


